# crontab
0 23 * * * /home/minecraft/mc-config/goto-bed.sh
0 12 * * * /home/minecraft/mc-config/goto-bed.sh